# target/linux/rockchip/armv8/base-files/etc/hotplug.d/net/01-usb-rename
#!/bin/sh

[[ "$ACTION" = add ]] || exit

case $(board_name) in
linkease,easepi-r2)
	if [[ ${DEVICENAME} != usb4g && ${DEVPATH} = /devices/platform/fc800000.usb/usb* ]]; then
		ip link set "${DEVICENAME}" name usb4g
	fi
	;;
esac
