# package/network/config/wifi-scripts/files/etc/hotplug.d/ieee80211/11-wifi-up-oneshot
#!/bin/sh

[ "${ACTION}" = "add" ] && {
	[ -f /tmp/.wifi-up ] && exit 0
	touch /tmp/.wifi-up

	/sbin/wifi isup || /sbin/wifi up
}
